# supervisord/hensei-web.ini
[program:hensei-web]
command=/opt/homebrew/bin/mise run dev
process_name=%(program_name)s
numprocs=1
directory=/Users/justin/Developer/Granblue/hensei-web
environment=HOME="/Users/justin",NODE_ENV="development",MISE_CONFIG_ROOT="/Users/justin/Developer/Granblue/hensei-web"
autostart=true
autorestart=unexpected
stopsignal=TERM
user=justin
stdout_logfile=/Users/justin/Developer/Granblue/hensei-web/logs/hensei-web.stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/Users/justin/Developer/Granblue/hensei-web/logs/hensei-web.stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
serverurl=AUTO
